# etc/systemd/system/kiosk-firstboot.service
[Unit]
Description=Kiosk first-boot one-time setup
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/env bash /usr/local/bin/kiosk-first-boot.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
